# gunicorn.conf
[program:gunicorn]
command=gunicorn -b localhost:80 -w 1 main:app
directory=/deploy/app