# gunicorn.conf
[program:gunicorn]
command=gunicorn -b localhost:5000 -w 1 main:app
directory=/deploy/app